# bin/kuromoji.service
[Unit]
Description=kuromoji
After=network.target

[Service]
User=root
Group=root
SyslogIdentifier=kuromoji.service
ExecStart=/root/kuromoji_api_docker/bin/run-dev.sh
WorkingDirectory=/root/kuromoji_api_docker/
Restart=always
RestartSec=5s
KillSignal=SIGINT
TimeoutStartSec=10
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target